# config/systemd/user/myarch-recover-internal-monitor.service
[Unit]
Description=Recover the internal monitor toggle when no external display is connected
Before=graphical-session-pre.target
ConditionPathExists=%h/.local/state/myarch/toggles/hypr/internal-monitor-disable.lua

[Service]
Type=oneshot
ExecStart=%h/.local/share/myarch/bin/myarch-hw-recover-internal-monitor

[Install]
WantedBy=graphical-session-pre.target
